# deploy/crowdtrans-web.service
[Unit]
Description=CrowdTrans Web Interface
After=network.target

[Service]
Type=simple
User=crowdtrans
Group=crowdtrans
WorkingDirectory=/opt/crowdtrans
ExecStart=/opt/crowdtrans/venv/bin/python -m crowdtrans.cli run-web
Restart=on-failure
RestartSec=10
EnvironmentFile=/opt/crowdtrans/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=crowdtrans-web

# Security
NoNewPrivileges=yes
ProtectHome=yes
ReadWritePaths=/opt/crowdtrans/data

# Allow binding to privileged ports (.env sets WEB_PORT=443).
# Without these, uvicorn fails with EACCES when started as the non-root
# crowdtrans user. NoNewPrivileges=yes still applies, since ambient caps
# are granted at exec time by systemd, not by the binary suid'ing up.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
